# templates/init-rhel.erb
#!/bin/sh -e
# bamboo startup script
#chkconfig: 2345 80 05
#description: bamboo
 
# Define some variables
# Name of app ( bamboo, Confluence, etc )
APP=bamboo
# Name of the user to run as
USER=<%= @user %>
# Location of application's bin directory
BASE=<%= @dir %>
# Location of Java JDK
export JAVA_HOME=<%= @java_home %>
 
case "$1" in
  # Start command
  start)
    echo "Starting $APP"
    /bin/su -m $USER -c "cd $BASE/logs && $BASE/bin/startup.sh &> /dev/null"
    ;;
  # Stop command
  stop)
    echo "Stopping $APP"
    /bin/su -m $USER -c "$BASE/bin/shutdown.sh &> /dev/null"
    echo "$APP stopped successfully"
    ;;
  # Restart command
  restart)
    $0 stop
    sleep 5
    $0 start
    ;;
  # Status command - there's no pid, this will do for now
  status)
    if ! netstat -alntp | grep 8085 | grep LISTEN | awk '{print $NF}' | grep java > /dev/null; then
      echo "$APP is not running."
      exit 3
    else
      echo "$APP is running"
    fi
    ;;
  *)
    echo "Usage: /etc/init.d/$APP {start|restart|stop}"
    exit 1
    ;;
esac
 
exit 0